# build/supervisord.conf.d/vncserver.conf
[program:vncserver]
command=vncserver :1 -geometry 1840x1000 -depth 24
stdout_logfile = /var/log/supervisor.%(program_name)s.log
stderr_logfile = /var/log/supervisor.%(program_name)s.log
autostart = true
autorestart = true
priority=15
