# build/supervisord.conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 8080 localhost:5900
stdout_logfile = /var/log/supervisor.%(program_name)s.log
stderr_logfile = /var/log/supervisor.%(program_name)s.log
autostart = true
autorestart = true
priority=20
